# rootfs/etc/services.d/viseron/run
#!/usr/bin/with-contenv bash

# Kill left over FFmpeg processes in case of a restart
pkill -f ffmpeg_

cd /src
exec justc-envdir /var/run/environment s6-setuidgid abc python3 -u -m viseron
