# rootfs/etc/services.d/viseron/run
#!/usr/bin/execlineb

justc-envdir /var/run/environment exec python3 -u /src/viseron/viseron.py
