# rootfs/etc/services.d/viseron/run
#!/usr/bin/execlineb

cd /src
justc-envdir /var/run/environment exec s6-setuidgid abc python3 -u -m viseron
