# rootfs/etc/services.d/viseron/run
#!/usr/bin/with-contenv bash

cd /src
exec justc-envdir /var/run/environment s6-setuidgid abc python3 -u -m viseron
