# rootfs/etc/services.d/postgres/run
#!/usr/bin/with-contenv bash

echo "Starting PostgreSQL Server..."
s6-setuidgid postgres /usr/lib/postgresql/12/bin/postgres -D /config/postgresql
